# deploy/hostinger/tot-dashboard-backup.timer
[Unit]
Description=Run the Tot Dashboard archive each hour

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=5m
AccuracySec=1m
Unit=tot-dashboard-backup.service

[Install]
WantedBy=timers.target
